# nginx/conf/project.cron
# -*- coding:utf-8 -*-
# author:phenix3443+github@gmail.com
# desc: 项目对应的crontab任务

# project crontab start
* * * * * <example-project-path>/nginx/sbin/monitor.sh
* * * * * <example-project-path>/nginx/sbin/report_falcon.sh

# project crontab end