# src/conf/project.cron
# -*- coding:utf-8 -*-
# author:phenix3443@gmail.com
# desc: 项目对应的crontab任务

# project crontab start

* * * * * nginx/sbin/report_falcon.sh

# project crontab end